# files/etc/uci-defaults/97-default-root-password
#!/bin/sh

if command -v chpasswd >/dev/null 2>&1; then
  echo 'root:root' | chpasswd && exit 0
fi

printf 'root\nroot\n' | passwd root >/dev/null 2>&1

exit 0
