# files/etc/uci-defaults/98-disable-mwan3
#!/bin/sh

if [ -x /etc/init.d/mwan3 ]; then
  /etc/init.d/mwan3 disable
fi

exit 0
